# systemd/x_start.service
[Unit]
Description=Monitor X server start

[Service]
Type=oneshot
ExecStart=

[Install]
WantedBy=default.target
